# Master/etc/udev/rules.d/90-ClusterNetworking.rules
SUBSYSTEM=="net", DRIVERS=="usb", RUN+="/opt/ClusterNetworkingAgent.sh"
